# fb-dashboard.service
[Unit]
Description=Frappe Manager Backup Dashboard
After=network.target

[Service]
Type=simple
User=frappe
WorkingDirectory=/home/frappe
# pipx installs to ~/.local/bin; ensure PATH includes it
Environment="PATH=/home/frappe/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/frappe/.local/bin/fb run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

